# tools/docker/daml-all-in-one/supervisord.conf
[supervisord]
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info

[program:daml]
# command=daml build && && daml sandbox --wall-clock-time --dar ./.daml/dist/quickstart-0.0.1.dar && daml json-api --config json-api-app.conf
command = bash -c 'daml build && daml sandbox --wall-clock-time --dar ./.daml/dist/quickstart-0.0.1.dar || { echo "Command failed" >&2; exit 1; }'
autostart=true
autorestart=true
stderr_logfile=/var/log/daml.err.log
stdout_logfile=/var/log/daml.out.log



[program:jsonapi]
# command=daml build && && daml sandbox --wall-clock-time --dar ./.daml/dist/quickstart-0.0.1.dar && daml json-api --config json-api-app.conf
command = daml json-api --config json-api-app.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/daml.err.log
stdout_logfile=/var/log/daml.out.log

[inet_http_server]
port = 0.0.0.0:9001



