# tools/docker/corda-all-in-one/supervisord.conf
[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[inet_http_server]
port = 0.0.0.0:9001

[program:sshd]
command=/usr/sbin/sshd -D -ddd
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.out.log
stderr_logfile=/var/log/sshd.err.log
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0

[program:dockerd]
command=dockerd-entrypoint.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:corda-a]
directory=/samples-kotlin/Advanced/obligation-cordapp/build/nodes/ParticipantA/
command=/run-party-a-node.sh
autostart=true
autorestart=false
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:corda-b]
directory=/samples-kotlin/Advanced/obligation-cordapp/build/nodes/ParticipantB
command=/run-party-b-node.sh
autostart=true
autorestart=false
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:corda-c]
directory=/samples-kotlin/Advanced/obligation-cordapp/build/nodes/ParticipantC
command=/run-party-c-node.sh
autostart=true
autorestart=false
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:corda-n]
directory=/samples-kotlin/Advanced/obligation-cordapp/build/nodes/Notary
command=/run-notary-node.sh
autostart=true
autorestart=false
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:corda-run-party-a-server]
directory=/samples-kotlin/Advanced/obligation-cordapp/
command=/run-party-a-server.sh
autostart=true
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:corda-run-party-b-server]
directory=/samples-kotlin/Advanced/obligation-cordapp/
command=/run-party-b-server.sh
autostart=true
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:corda-run-party-c-server]
directory=/samples-kotlin/Advanced/obligation-cordapp/
command=/run-party-c-server.sh
autostart=true
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
