# tools/all-in-one/corda/supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]

[inet_http_server]
port = 127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:sshd]
command=/usr/sbin/sshd -D

[program:partyA]
command=java -jar /opt/corda/partyA/corda.jar --base-directory=/opt/corda/partyA
stderr_logfile=/var/log/partyA.err.log
stdout_logfile=/var/log/partyA.out.log

[program:partyB]
command=java -jar /opt/corda/partyB/corda.jar --base-directory=/opt/corda/partyB
stderr_logfile=/var/log/partyB.err.log
stdout_logfile=/var/log/partyB.out.log
