# tools/docker/indy-testnet/validator/supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
logLevel = debug
pidfile = /tmp/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
user = indy
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[program:validator]
command=/bin/bash -c 'echo start-server && TEST_POOL_IP=172.16.0.2 python3 main.py'
directory=/home/indy/cactus_validator_socketio/validator-python
