# packages/cactus-plugin-satp-hermes/supervisord.conf
[supervisord]
logfile=/opt/cacti/satp-hermes/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
user=root

[program:satp-gateway]
command=/usr/local/bin/node index.js
autostart=true
stderr_logfile=/opt/cacti/satp-hermes/logs/satp-gateway-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile=/opt/cacti/satp-hermes/logs/satp-gateway-output.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
