# packages/cactus-plugin-satp-hermes/supervisord.conf
[supervisord]
logfile=/opt/cacti/satp-hermes/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
user=root

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:dockerd]
command=/usr/local/bin/dockerd
autostart=true
autorestart=true
startretries=10
priority=100

[program:satp-gateway]
command=/bin/bash -c "/usr/local/bin/node index.js \
    > >(/opt/cacti/satp-hermes/gateway_log_controller.sh stdout) \
    2> >(/opt/cacti/satp-hermes/gateway_log_controller.sh stderr)"
autostart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=999
