# tools/docker/fabric-all-in-one/supervisord.conf
[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info


[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/sshd.err.log
stdout_logfile=/var/log/sshd.out.log

[program:dockerd]
command=dockerd-entrypoint.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/dockerd.err.log
stdout_logfile=/var/log/dockerd.out.log

[program:fabric-network]
command=/run-fabric-network.sh
autostart=true
autorestart=unexpected
stderr_logfile=/var/log/fabric-network.err.log
stdout_logfile=/var/log/fabric-network.out.log

[inet_http_server]
port = 0.0.0.0:9001

