# tools/docker/fabric-all-in-one/supervisord.conf
[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info

# Commented CA for future use, if needed
#[program:ca]
#command=/etc/hyperledger/fabric-ca-server/start_ca.sh
#autostart=true
#autorestart=true
#stderr_logfile=/var/log/fabric-ca.err.log
#stdout_logfile=/var/log/fabric-ca.out.log

# FXIME: We need DIND so that the AIO container is truly all-in-one
# [program:dind]
# command=/usr/local/bin/dind
# stderr_logfile=/var/log/dind.err.log
# stdout_logfile=/var/log/dind.out.log

[program:sshd]
command=/usr/sbin/sshd -D
stderr_logfile=/var/log/sshd.err.log
stdout_logfile=/var/log/sshd.out.log

[program:orderer]
command=/etc/hyperledger/fabric/orderer/start_orderer.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/orderer.err.log
stdout_logfile=/var/log/orderer.out.log

[program:peer]
command=/etc/hyperledger/fabric/peer/start_peer.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/peer.err.log
stdout_logfile=/var/log/peer.out.log

[program:channel]
command=/etc/hyperledger/fabric/peer/join_channel.sh
autostart=true
autorestart=unexpected
startsecs=0
exitcodes=0
startretries=2
stderr_logfile=/var/log/channel.err.log
stdout_logfile=/var/log/channel.out.log

[inet_http_server]
port = 0.0.0.0:9001
